# roles/monitoring/templates/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
ExecStart=/usr/share/prometheus/node_exporter

[Install]
WantedBy=default.target
